# etc/profile.d/sunray.sh
PATH=$PATH:/opt/SUNWut/bin:/opt/SUNWuttsc/bin:/opt/SUNWut/sbin
MANPATH=$MANPATH:/opt/SUNWut/share/man:/opt/SUNWuttsc/share/man:/opt/SUNWkio/share/man
export PATH MANPATH
